# docker/quadlet/dealdesk-db.container
[Unit]
Description=PostgreSQL for DealDesk

[Container]
Image=docker.io/library/postgres:17-alpine
ContainerName=dealdesk-db
Pod=dealdesk.pod
Volume=dealdesk-pgdata:/var/lib/postgresql/data
EnvironmentFile=%h/.config/containers/systemd/dealdesk.env
HealthCmd=pg_isready -U $POSTGRES_USER -d $POSTGRES_DB -h localhost || exit 1
HealthInterval=15s
HealthTimeout=5s
HealthRetries=5

[Service]
Restart=on-failure
TimeoutStartSec=60

[Install]
WantedBy=default.target
